# hf-space/supervisord.conf
# ──────────────────────────────────────────────────────────────
# SmartLearn AI - supervisord 进程编排配置
# ──────────────────────────────────────────────────────────────
# 在 HF Spaces 单容器内管理 3 个进程:
#   1. nginx         :7860  反向代理 + 静态托管（唯一对外端口）
#   2. uvicorn api   :8000  FastAPI 业务后端
#   3. uvicorn ai-engine :8001  AI 引擎（RAG/游戏/对话）
# ──────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

# ── nginx 反向代理 ──
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ── API 后端 (FastAPI) ──
# workers=1 节省内存,加快冷启动;proxy_headers 让后端识别真实 IP/协议
[program:api]
directory=/app
command=uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1 --proxy-headers
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

# ── AI 引擎 (FastAPI) ──
# workers=1 节省内存;ai-engine 启动时加载知识库,需更长 startsecs
[program:ai-engine]
directory=/app/ai-engine
command=uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 1 --proxy-headers
autostart=true
autorestart=true
startsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
